# etc/systemd/specter.service
[Unit]
Description=Specter — prove the absence of change
After=network.target

[Service]
Type=exec
ExecStart=@SPECTER_BIN@ run --config @SPECTER_CONF@
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

User=@SPECTER_USER@
Group=@SPECTER_GROUP@

StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true

ReadOnlyPaths=/srv
ReadWritePaths=/var/log/specter.log

[Install]
WantedBy=multi-user.target
